# ODOO-POSTGRESSQL/roles/borg-backup-client/files/borg-backup.timer
[Unit]
Description=Borg Backup
[Timer]
OnUnitActiveSec=60min
[Install]
WantedBy=timers.target